# deploy/tw.observe.mayor2026.pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>tw.observe.mayor2026.pipeline</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>-lc</string>
      <string>cd '/Users/REPLACE_ME/mayor2026' &amp;&amp; /usr/bin/python3 scripts/run_pipeline.py --publish-pages</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/REPLACE_ME/mayor2026</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>MAYOR_OBSERVE_HOME</key>
      <string>/Users/REPLACE_ME/mayor2026</string>
    </dict>

    <!-- MAYOR_RSSHUB_BASE is intentionally left out here: when this deploys
         beside a machine that already runs RSSHub locally (e.g. the same
         host as Harmonica-in-Taiwan), set MAYOR_RSSHUB_BASE=http://127.0.0.1:1200
         in .env instead of here, so it's easy to flip back to the public
         https://rss.observe.tw default without editing this plist. -->

    <key>StartCalendarInterval</key>
    <array>
      <dict><key>Hour</key><integer>0</integer><key>Minute</key><integer>0</integer></dict>
      <dict><key>Hour</key><integer>6</integer><key>Minute</key><integer>0</integer></dict>
      <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>0</integer></dict>
      <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/REPLACE_ME/mayor2026/logs/pipeline.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/REPLACE_ME/mayor2026/logs/pipeline.err.log</string>
  </dict>
</plist>
