# deploy/tw.observe.mayor2026.pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>tw.observe.mayor2026.pipeline</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>-lc</string>
      <string>cd '/Users/REPLACE_ME/mayor2026' &amp;&amp; /usr/bin/python3 scripts/run_pipeline.py --publish-pages</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/REPLACE_ME/mayor2026</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>MAYOR_OBSERVE_HOME</key>
      <string>/Users/REPLACE_ME/mayor2026</string>
      <key>MAYOR_RSSHUB_BASE</key>
      <string>https://rss.observe.tw</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <!-- Every 30 minutes, same cadence as Harmonica-in-Taiwan; the pipeline
         lock and per-source fetch intervals (e.g. Instagram 12h) keep this
         from hammering any single platform. -->
    <key>StartInterval</key>
    <integer>1800</integer>

    <key>StandardOutPath</key>
    <string>/Users/REPLACE_ME/mayor2026/logs/pipeline.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/REPLACE_ME/mayor2026/logs/pipeline.err.log</string>
  </dict>
</plist>
